# service/lambada-web.service
[Unit]
Description=Lambada Web (scan listing server)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/workspace/lambada
ExecStart=/home/pi/workspace/lambada/lambada-web
Restart=always
# Backstop fd ceiling -- see docs/COWORK.md (issue #2) for why.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
