# service/lambada-mta.service
[Unit]
Description=Lambada MTA (SMTP scan exporter)
After=network.target

[Service]
Type=simple
User=woodie
WorkingDirectory=/home/woodie/workspace/lambada
ExecStart=/home/woodie/workspace/lambada/lambada-mta
Restart=always
# Matching lambada-web.service's LimitNOFILE -- mta hasn't shown the same
# symptom, but there's no reason to leave it on a stingier default fd
# ceiling than its sibling service.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
